# data/xdg-desktop-portal-tau.service.in
[Unit]
Description=Portal service (tauOS implementation)
After=graphical-session.target
Requisite=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.tau
ExecStart=@libexecdir@/xdg-desktop-portal-tau