# data/xdg-desktop-portal-tau.service.in
[Unit]
Description=Portal service (tauOS implementation)

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.tau
ExecStart=@libexecdir@/xdg-desktop-portal-tau
