# data/xdg-desktop-portal-gnome.service.in
[Unit]
Description=Portal service (gnomeOS implementation)

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.gnome
ExecStart=@libexecdir@/xdg-desktop-portal-gnome
